# packages/django-channel-tasks/django_channel_tasks-0.4.3.dev4-py3-none-any.whl/django_channel_tasks-0.4.3.dev4.data/data/conf/nginx-unit/channel-tasks.template.service
[Unit]
Description=channel-tasks

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=${CHANNEL_TASKS_BIN_DIR}/restart-channel-tasks.sh
Environment="CHANNEL_TASKS_DOCKER_HOME=${CHANNEL_TASKS_DOCKER_HOME}"

[Install]
WantedBy=multi-user.target
